# async_method/systemd/proxy-pool-scanner-restart.service
[Unit]
Description=Restart proxy-pool-scanner.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart proxy-pool-scanner.service
