# async_method/systemd/proxy-pool-checker-restart.service
[Unit]
Description=Restart proxy-pool-checker.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart proxy-pool-checker.service
