# async_method/systemd/proxy-pool-checker-restart.timer
[Timer]
OnActiveSec=2
OnUnitActiveSec=12

[Install]
WantedBy=timer.target
