# systemd/proxy-pool-sort-restart.timer
[Timer]
OnActiveSec=5
OnUnitActiveSec=5
#Persistent=true

[Install]
WantedBy=timer.target
