# systemd/proxy-pool-sort-restart.service
[Unit]
Description=Restart proxy-pool-sort.service

[Service]
Type=oneshot
ExecStart=/bin/systemctl --user restart proxy-pool-sort.service
