# systemd/proxy-pool-scanner-restart.timer
[Timer]
OnActiveSec=180
OnUnitActiveSec=180

[Install]
WantedBy=timer.target
