# systemd/proxy-pool-scanner-restart.timer
[Timer]
OnActiveSec=300
OnUnitActiveSec=300

[Install]
WantedBy=timer.target
